# debian/nv-opof.postinst
systemctl daemon-reload
systemctl enable nv-opof
